# ynh-package/conf/systemd.service
[Unit]
Description=Nexora Control Plane
After=network.target

[Service]
Type=simple
User=__APP__
WorkingDirectory=__INSTALL_DIR__
Environment=NEXORA_STATE_PATH=__DATA_DIR__/state.json
Environment=NEXORA_API_TOKEN_FILE=__DATA_DIR__/api-token
Environment=NEXORA_API_TOKEN_ROLE_FILE=/etc/nexora/api-token-roles.json
Environment=NEXORA_OPERATOR_ONLY_ENFORCE=1
Environment=NEXORA_DEPLOYMENT_SCOPE=operator
Environment=NEXORA_PROFILE=__NEXORA_PROFILE__
Environment=NEXORA_DEPLOYMENT_MODE=__DEPLOYMENT_MODE__
Environment=NEXORA_ENROLLMENT_MODE=__ENROLLMENT_MODE__
Environment=NEXORA_TENANT_DEFAULT_TIER=__TENANT_DEFAULT_TIER__
Environment=NEXORA_SAAS_OPERATOR_URL=__SAAS_OPERATOR_URL__
Environment=NEXORA_CONTROL_PLANE_HOST=127.0.0.1
Environment=NEXORA_CONTROL_PLANE_PORT=__PORT__
ExecStart=__INSTALL_DIR__/venv/bin/nexora-control-plane
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=__DATA_DIR__ /tmp/nexora-export /var/log/yunohost-mcp-server
ProtectHome=false

[Install]
WantedBy=multi-user.target
